# thinkpadbacklight-restart.timer
[Unit]
Description=Run thinkpadbacklight restart hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=multi-user.target
