# systemd/thinkpadbacklight.service
[Unit]
Description=Thinkpad keyboard backlight

[Service]
ExecStart=/bin/bash /home/theo/.scripts/tpbacklight
[Install]
WantedBy=default.target
